# install/obgs-farmer.service
[Unit]
Description=Open Banner Grabbing System

[Service]
Type=simple
ExecStart=/srv/obgs/farmer/farmer.py

[Install]
WantedBy=multi-user.target
