# install/obgs-api.service
[Unit]
Description=Open Banner Grabbing System

[Service]
Type=simple
ExecStart=/srv/obgs/api/server.py

[Install]
WantedBy=multi-user.target
